# SOURCES_TAR/rgm-base/backup/rgmbackup.cron
15 3 * * * root /srv/rgm/backup/bin/rgm-backup.sh &> /dev/null
